# docker/supervisord.conf
# Dropped into /etc/supervisor/conf.d/ — loaded via [include] from the main supervisord.conf.
# Only [program:] sections belong here; global [supervisord] settings live in the main conf.

# ---------------------------------------------------------------------------
# nginx — edge server; serves SPA and proxies /api/* to uvicorn
# Stays up while uvicorn restarts so the SPA shell remains reachable.
# ---------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


# ---------------------------------------------------------------------------
# uvicorn — FastAPI backend; single worker mandatory (DESIGN §19.2)
# Restarted by supervisor with retries; nginx stays up during restarts.
# stopwaitsecs is generous so in-flight saves can flush (DESIGN §19.7).
#
# stdout/stderr go to the container's own stdout/stderr (not a file inside the
# container) so the backend's logs — including the import/startup traceback
# behind an `exit status 1` crash-loop — reach `docker logs` and any hosting
# platform's log aggregator. A file sink is invisible in deployments where you
# cannot exec into the container. maxbytes=0 disables rotation, which is
# required when the sink is a non-seekable device rather than a regular file.
# ---------------------------------------------------------------------------
[program:uvicorn]
command=python -m uvicorn backend.app.main:app --host 127.0.0.1 --port 8000 --workers 1 --log-level info --no-access-log
directory=/app
autostart=true
autorestart=true
startretries=10
startsecs=5
stopsignal=SIGTERM
stopwaitsecs=120
priority=20
environment=PYTHONPATH="/app"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
